# systemd/llm-gateway.service
[Unit]
Description=Ollama LLM Gateway API
After=network.target ollama.service

[Service]
Type=simple
WorkingDirectory=/opt/llm-gateway/server
ExecStart=/opt/llm-gateway/venv/bin/python gateway.py
Restart=always
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target